# cfddns.service
[Unit]
Description="Cloudflare DDNS updater"

[Service]
#edit me 
ExecStart=/path/to/cfddns.sh
